# docker/service/avahi-daemon/run
#!/bin/sh
s6-svwait ../dbus-daemon
exec avahi-daemon
