# docker/service/syndicate-server/run
#!/bin/sh
s6-svwait ../avahi-daemon
cd /data/packages/server
exec node lib/index.js
